# patches/etc/systemd/system/ocs-sample.service
[Unit]
Description=OCS sample service

[Service]
ExecStart=/etc/ocs/hello-world
Type=oneshot

[Install]
WantedBy=multi-user.target
